# Client-Server sockets/MyDevServer.service
[Unit]
Description=Server

[Service]
ExecStart=/home/user/Practice/day9/Server/build/server

[Install]
WantedBy=multi-user.target
